# hs-proxy.service
[Unit]
Description = ssl proxy

[Service]
ExecStart = HOME/ssh-proxy-start.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
